# container_scripts/services/app/run
#!/bin/bash -e

cd /home/app/app
/sbin/setuser app npm install
/sbin/setuser app node app.js
